# package/base-files/files/etc/uci-defaults/70_pptp_l2tp_mover
#!/bin/sh

. /lib/functions.sh

cp /rom/etc/ppp/options.pptpd /etc/ppp/options.pptpd
cp /rom/etc/ppp/options.xl2tpd /etc/ppp/options.xl2tpd
exit 0
